# recipes-core/swupdate/swupdate/swupdate.service
[Unit]
Description=swupdate
ConditionPathExists=!/data

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
ExecStart=/usr/bin/swupdate.sh

[Install]
WantedBy=multi-user.target
